# deploy/etc/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; (the path to the socket file)

[supervisord]
logfile=/tmp/supervisord.log ; (main log file;default $CWD/supervisord.log)
loglevel=info               ; (log level;default info; others: debug,warn,trace)
pidfile=/tmp/supervisord.pid ; (supervisord pidfile;default supervisord.pid)
nodaemon=false              ; (start in foreground if true;default false)
minfds=1024                 ; (min. avail startup file descriptors;default 1024)
minprocs=200                ; (min. avail process descriptors;default 200)
environment=DJANGO_ENVIRON={env.project_env}

; the below section must remain in the config file for RPC
; (supervisorctl/web interface) to work, additional interfaces may be
; added by defining them in separate rpcinterface: sections
[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use a unix:// URL  for a unix socket


[program:gunicorn]
command={env.venv_root}bin/python {env.src_root}manage.py run_gunicorn --preload --workers=3 --log-file=/tmp/gunicorn.log --backlog=2048 --pid=/tmp/gunicorn.pid --bind=127.0.0.1:8000 --timeout=300 --settings=settings.prod
; user=nobody
autostart=true
autorestart=true
redirect_stderr=true

;[program:solr]
;directory={solr.start_root}
;command=java -jar start.jar
;startsecs=20
;stopsignal=QUIT
;;user=
;autostart=true
;autorestart=true
;redirect_stderr=true
;stdout_logfile=/tmp/supervisord.solr.log


; [program:celery]
; command={env.venv_root}bin/python {env.src_root}manage.py celeryd -f /tmp/celeryd.log -l INFO -v 0 --settings=settings.prod


; The [include] section can just contain the "files" setting.  This
; setting can list multiple files (separated by whitespace or
; newlines).  It can also contain wildcards.  The filenames are
; interpreted as relative to this file.  Included files *cannot*
; include files themselves.

;[include]
;files = relative/directory/*.ini
